# systemd/liturgical.timer
[Unit]
Description=Daily Liturgical Image Update

[Timer]
OnCalendar=*-*-* 00:01:00
Persistent=true

[Install]
WantedBy=timers.target 