# systemd/scriptura-api.service
[Unit]
Description=Scriptura API - Local Bible API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{PROJECT_DIR}}/scriptura-api
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -c "source venv/bin/activate && uvicorn main:app --host 0.0.0.0 --port 8081"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=false
ProtectHome=false
ReadWritePaths={{PROJECT_DIR}}/scriptura-api

# Resource limits for low-memory systems
MemoryLimit=128M
MemoryAccounting=true
MemoryHigh=96M
MemoryMax=128M
LimitNOFILE=1024
LimitNPROC=50

# Additional resource constraints
LimitAS=128M
LimitRSS=128M
LimitDATA=128M
LimitFSIZE=10M
LimitSTACK=8M
LimitCORE=0
LimitMEMLOCK=64K
LimitLOCKS=32
LimitSIGPENDING=32
LimitMSGQUEUE=819200
LimitNICE=0
LimitRTPRIO=0
LimitRTTIME=1000000

[Install]
WantedBy=multi-user.target
